# com.aiautotool.aiagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aiautotool.aiagent</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/vkct/aiagent/ai_api.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>WorkingDirectory</key>
    <string>/Users/vkct/aiagent</string>
    <key>StandardOutPath</key>
    <string>/Users/vkct/aiagent/ai_api.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/vkct/aiagent/ai_api.log</string>
</dict>
</plist>
